# templates/default/vertx_runlevel.erb
#!/bin/bash
#
# vertx Startup script for Vert.x
#

<% if node[:platform_family] == "rhel" -%>
# chkconfig: 2345 85 15
# description: vertx
#
### BEGIN INIT INFO
# Provides: $vertx
### END INIT INFO

if [ -f /etc/sysconfig/vertx ]; then
  . /etc/sysconfig/vertx_deploy
fi

<% elsif node[:platform_family] == "debian" -%>
### BEGIN INIT INFO
# Provides:             vertx
# Required-Start:       $remote_fs $syslog
# Required-Stop:        $remote_fs $syslog
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    Vert.x server
### END INIT INFO

if [ -f /etc/default/vertx_deploy ]; then
  . /etc/default/vertx_deploy
fi
<% end -%>

VERTX_USER=${VERTX_USER:-"vertx"}

if [ "$EUID" = "0" ]; then
  SUBIT="su - $VERTX_USER -c "
else
  if [ "$USER" = "vertx" ]; then
     SUBIT=""
  else
     echo "Please switch to user vertx !"
     exit 1
  fi
fi

STATUS_FILE=<%= node[:status_file] %>
PID_FILE=<%= node[:pid_file] %>

case "$1" in
     start)
        echo "Starting Vert.x <%= node[:version] %>"
        if [ -z "$VERTX_MODULE" ]; then
          exit 0
        fi
        <% if node[:platform_family] == "debian" -%>
        START_SCRIPT="authbind <%= node[:dir] %>/vertx/bin/vertx runzip <%= node[:dir] %>/deploy/packages/$VERTX_MODULE"
        <% else %>
        START_SCRIPT="<%= node[:dir] %>/vertx/bin/vertx runzip <%= node[:dir] %>/deploy/packages/$VERTX_MODULE"
        <% end %>
        if [ -z "$SUBIT" ]; then
          $START_SCRIPT > <%= node[:log_dir] %>/console.log &
          echo ONLINE > $STATUS_FILE
        else
          touch <%= node[:log_dir] %>/console.log
          chown vertx:vertx <%= node[:log_dir] %>/console.log
          $SUBIT "$START_SCRIPT" > <%= node[:log_dir] %> &
          $SUBIT "echo ONLINE > $STATUS_FILE"
        fi
        echo `ps awwx | grep "vertx.java.platform" | grep -v grep | awk '{print $1}'` > $PID_FILE
    ;;
    stop)
        echo "Stopping Vert.x <%= node[:version] %>"
        if [ -z "$VERTX_MODULE" ]; then
          exit 0
        fi
        PID=`ps awwx | grep $VERTX_MODULE | grep -v grep | grep java | awk '{ print $1 }'`
        if [ -n "$PID" ]; then
          STOP_SCRIPT="kill $PID"
          if [ -z "$SUBIT" ]; then
            echo OFFLINE > $STATUS_FILE
            $STOP_SCRIPT > <%= node[:log_dir] %>/console.log &
          else
            touch <%= node[:log_dir] %>/console.log
            chown vertx:vertx <%= node[:log_dir] %>/console.log
            $SUBIT "echo OFFLINE > $STATUS_FILE"
            $SUBIT "$STOP_SCRIPT" > <%= node[:log_dir] %>/console.log &
          fi
        fi
        if [ -f $PID_FILE ]; then
          rm -f $PID_FILE
        fi
    ;;
    status)
        echo "Status of Vert.x <%= node[:version] %>"
        STATUS_SCRIPT="<%= node[:dir] %>/vertx/bin/vertx_status"
        if [ -z "$SUBIT" ]; then
            $STATUS_SCRIPT
            RETVAL=$?
        else
            $SUBIT "$STATUS_SCRIPT"
            RETVAL=$?
        fi
  
        exit $RETVAL
    ;;
    restart)
        echo "Stopping Vert.x <%= node[:version] %>"
        if [ -z "$VERTX_MODULE" ]; then
          exit 0
        fi
        PID=`ps awwx | grep $VERTX_MODULE | grep -v grep | grep java | awk '{ print $1 }'`
        if [ -n "$PID" ]; then
          STOP_SCRIPT="kill $PID"
          if [ -z "$SUBIT" ]; then
            echo OFFLINE > $STATUS_FILE
            $STOP_SCRIPT > <%= node[:log_dir] %>/console.log &
          else
            touch <%= node[:log_dir] %>/console.log
            chown vertx:vertx <%= node[:log_dir] %>/console.log
            $SUBIT "echo OFFLINE > $STATUS_FILE"
            $SUBIT "$STOP_SCRIPT" > <%= node[:log_dir] %>/console.log &
          fi
        fi

        sleep 5

        echo "Starting Vert.x <%= node[:version] %>"
        if [ -z "$VERTX_MODULE" ]; then
          exit 0
        fi
        <% if node[:platform_family] == "debian" -%>
        START_SCRIPT="authbind <%= node[:dir] %>/vertx/bin/vertx runzip <%= node[:dir] %>/deploy/packages/$VERTX_MODULE"
        <% else %>
        START_SCRIPT="<%= node[:dir] %>/vertx/bin/vertx runzip <%= node[:dir] %>/deploy/packages/$VERTX_MODULE"
        <% end %>
        if [ -z "$SUBIT" ]; then
          $START_SCRIPT > <%= node[:log_dir] %>/console.log &
          echo ONLINE > $STATUS_FILE
        else
          touch <%= node[:log_dir] %>/console.log
          chown vertx:vertx <%= node[:log_dir] %>/console.log
          $SUBIT "echo OFFLINE > $STATUS_FILE"
          $SUBIT "$STOP_SCRIPT" > <%= node[:log_dir] %>/console.log &
        fi
    ;;
    *)
        echo "Usage: /etc/init.d/vertx {start|stop|restart}"
        exit 1
    ;;
esac

exit 0
